# scripts/vera-gateway-update.path
[Unit]
Description=Watch for Vera Gateway update requests

[Path]
PathExists=/var/lib/vera-updater/requests/request.json
Unit=vera-gateway-update.service

[Install]
WantedBy=multi-user.target
